# .crontab
# src: https://github.com/paulirish/dotfiles/blob/master/crontab
#
# Usage:
#  - crontab -e
#  - crontab $HOME/.crontab (to notify cron of new changes)
#  - chmod +x .cron/*.sh if new scripts were added

#=====================================================================#
#     0-23   0-59               0-30  0-6   1-12                      #
#     |_    _|                 _|    _|    _|                         #
#       |  |                  |     |     |            * any          #
#    At HH:MM on day-of-month MD on WD in M            - range        #
#       |  |                  |     |     |            , array        #
#     __|__|            ______|     |   __|            / step         #
#    |  |            __|         ___|__|                              #
#   _|  |__         |         __|   |__                               #
#  |       |        |        |         |                              #
#  |       |        |        |         |                              #
#=====================================================================#
#  Min  |  Hour  |  Mday  |  Month  |  Wday  |  Command               #
#=====================================================================#
   0       3        *        *         *        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/daily-night" > /dev/null 2>&1
   */5     *        *        *         *        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/5min-any" > /dev/null 2>&1
   0       *        *        *         *        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/hourly-any" > /dev/null 2>&1
   20      17       *        *         1-5      sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/weekday-afternoon" > /dev/null 2>&1
   0       17       *        *         0        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/daily-afternoon" > /dev/null 2>&1
   0       7        *        *         0        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/daily-morning" > /dev/null 2>&1
   0       0        *        *         5        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/weekend-any" > /dev/null 2>&1
   0       0        1        *         *        sh $HOME/.bin/cron-run-dir.sh "$HOME/.cron/monthly-any" > /dev/null 2>&1
#=====================================================================#
